# rlmtools.cron
0,30 * * * * apache bash /afs/unity/web/l/linux/secure/xmlrpc/rrd-update.sh
43   2 * * * apache bash /afs/unity/web/l/linux/secure/xmlrpc/dbcron.sh

