# rlmtools.cron
0,30 * * * * apache python /home/slack/projects/xmlrpc/rrd-update.sh
43   2 * * * apache python /home/slack/projects/xmlrpc/dbcron.sh

